# packaging/rpm/.build/secureusb-1.0.0/packaging/debian/postinst
#!/bin/bash

set -e

CONFIG_DIR="/var/lib/secureusb"
POINTER_DIR="/etc/secureusb"
POINTER_FILE="$POINTER_DIR/config_dir"

echo "Configuring SecureUSB..."
python3 -m pip install --break-system-packages -r /opt/secureusb/requirements.txt

mkdir -p "$CONFIG_DIR"
chmod 700 "$CONFIG_DIR"
mkdir -p "$POINTER_DIR"
echo "$CONFIG_DIR" > "$POINTER_FILE"

systemctl daemon-reload
udevadm control --reload-rules

# Update icon cache
if command -v gtk-update-icon-cache &> /dev/null; then
    gtk-update-icon-cache -f -t /usr/share/icons/hicolor 2>/dev/null || true
fi

systemctl enable secureusb.service || true
systemctl restart secureusb.service || true

echo "SecureUSB installed. Run 'secureusb-setup' to configure TOTP."
