# packaging/rpm/.build/secureusb-1.0.0/packaging/debian/prerm
#!/bin/bash

set -e

systemctl stop secureusb.service 2>/dev/null || true
systemctl disable secureusb.service 2>/dev/null || true
