# packaging/debian/postinst
#!/bin/bash

set -e

CONFIG_DIR="/var/lib/secureusb"
POINTER_DIR="/etc/secureusb"
POINTER_FILE="$POINTER_DIR/config_dir"

echo "Configuring SecureUSB..."
python3 -m pip install --break-system-packages -r /opt/secureusb/requirements.txt

mkdir -p "$CONFIG_DIR"
chmod 700 "$CONFIG_DIR"
mkdir -p "$POINTER_DIR"
echo "$CONFIG_DIR" > "$POINTER_FILE"

systemctl daemon-reload
udevadm control --reload-rules

# Update icon cache
if command -v gtk-update-icon-cache &> /dev/null; then
    gtk-update-icon-cache -f -t /usr/share/icons/hicolor 2>/dev/null || true
fi

# Try to install and enable GNOME AppIndicator extension for GNOME users
if command -v gnome-shell &> /dev/null; then
    echo "Detected GNOME Shell - ensuring AppIndicator extension is available"

    # Try to install the extension package if not already installed
    if ! dpkg -l | grep -q gnome-shell-extension-appindicator; then
        echo "Installing GNOME AppIndicator extension..."
        apt-get install -y gnome-shell-extension-appindicator 2>/dev/null || true
    fi

    # Find the actual user (not root) to enable the extension for
    ACTUAL_USER="${SUDO_USER:-$USER}"
    ACTUAL_HOME=$(getent passwd "$ACTUAL_USER" | cut -d: -f6)

    # Try to enable the extension - check for different extension IDs
    for EXT_ID in "ubuntu-appindicators@ubuntu.com" "appindicatorsupport@rgcjonas.gmail.com"; do
        if sudo -u "$ACTUAL_USER" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$(id -u "$ACTUAL_USER")/bus" \
           gnome-extensions list --user 2>/dev/null | grep -q "$EXT_ID"; then
            echo "Enabling AppIndicator extension: $EXT_ID"
            sudo -u "$ACTUAL_USER" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$(id -u "$ACTUAL_USER")/bus" \
                gnome-extensions enable "$EXT_ID" 2>/dev/null || true
            break
        fi
    done
fi

systemctl enable secureusb.service || true
systemctl restart secureusb.service || true

echo "SecureUSB installed. Run 'secureusb-setup' to configure TOTP."
