# macos/org.secureusb.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.secureusb.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/opt/secureusb/src/daemon/service.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/secureusb</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Library/Logs/SecureUSB/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Logs/SecureUSB/daemon.error.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>Disabled</key>
    <false/>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
